# autostart/.config/autostart/gnome-terminal.desktop
# - Notes
#   - `--geometry=...`: window size
#
[Desktop Entry]
Type=Application
Exec=gnome-terminal --geometry=135x68+0
